# Linux/systemd/user/rclone-googledrive.timer
[Unit]
Description=Run Rclone Backup every 2 hours

[Timer]
OnBootSec=30min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
